# docker/crontab
* * * * * php /fleetbase/api/artisan schedule:run
