# docker/crontab
* * * * * php /var/www/html/api/artisan schedule:run
